# systemd/.config/systemd/user/syncthing.service.d/override.conf
[Service]
Nice=6
IOSchedulingClass=2
IOSchedulingPriority=4
